# repo/docker/crontab
# Appointment auto-expire: run every 15 minutes
*/15 * * * * cd /var/www/html && php think appointment:expire >> /var/log/cron.log 2>&1

# Nightly risk scoring: run at 2:00 AM
0 2 * * * cd /var/www/html && php think risk:score >> /var/log/cron.log 2>&1

# Monthly audit log archive: run 1st of month at 3:00 AM
0 3 1 * * cd /var/www/html && php think audit:archive >> /var/log/cron.log 2>&1
